# packages/luci-app-mt5700m/root/etc/uci-defaults/90-mt5700m
#!/bin/sh

uci -q get mt5700m.settings >/dev/null 2>&1 || uci -q set mt5700m.settings='mt5700m'
uci -q set mt5700m.settings.enabled="${MT5700M_ENABLED:-1}"
uci -q set mt5700m.settings.mode="$(uci -q get mt5700m.settings.mode || echo auto)"
uci -q set mt5700m.settings.at_port="$(uci -q get mt5700m.settings.at_port || echo '')"
uci -q set mt5700m.settings.host="$(uci -q get mt5700m.settings.host || echo 192.168.8.1)"
uci -q set mt5700m.settings.port="$(uci -q get mt5700m.settings.port || echo 20249)"
uci -q set mt5700m.settings.timeout="$(uci -q get mt5700m.settings.timeout || echo 8)"
uci -q commit mt5700m

exit 0
