# packages/luci-app-h5000m-fancontrol/root/etc/uci-defaults/90-h5000m-fancontrol
#!/bin/sh

/etc/init.d/h5000m-fancontrol enable >/dev/null 2>&1 || true
/etc/init.d/h5000m-fancontrol restart >/dev/null 2>&1 || true

exit 0
